# tx.socket
[Unit]
Description=Socket to communicate with myservice

[Socket]
ListenFIFO=/tmp/modem.socket
